# server/vps-agent.service
[Unit]
Description=VPS Manager Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-agent
ExecStart=/opt/vps-agent/vps-agent -c /opt/vps-agent/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-agent

# Security hardening (comment out if issues arise)
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
